# infra/systemd/edge-policy-bridge-mqtt.service
[Unit]
Description=Edge Policy Hub - MQTT Bridge Service
Documentation=https://github.com/acailic/edge-policy-hub
After=network-online.target edge-policy-enforcer.service
Requires=edge-policy-enforcer.service
PartOf=edge-policy-hub.target

[Service]
Type=simple
ExecStart=/usr/local/bin/edge-policy-bridge-mqtt
Restart=on-failure
RestartSec=5s
User=edge-policy
Group=edge-policy
WorkingDirectory=/var/lib/edge-policy-hub
Environment="MQTT_HOST=0.0.0.0"
Environment="MQTT_PORT=1883"
Environment="MQTT_BROKER_NAME=edge-policy-mqtt"
Environment="ENFORCER_URL=http://127.0.0.1:8181"
Environment="TOPIC_NAMESPACE_PATTERN={tenant_id}/#"
Environment="LOG_LEVEL=info"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=edge-policy-bridge-mqtt
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/edge-policy-hub
CapabilityBoundingSet=
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
Alias=bridge-mqtt.service
