# infra/systemd/edge-policy-enforcer.service
[Unit]
Description=Edge Policy Hub - OPA Enforcer Service
Documentation=https://github.com/acailic/edge-policy-hub
After=network-online.target
Wants=network-online.target
PartOf=edge-policy-hub.target

[Service]
Type=simple
ExecStart=/usr/local/bin/edge-policy-enforcer
Restart=on-failure
RestartSec=5s
User=edge-policy
Group=edge-policy
WorkingDirectory=/var/lib/edge-policy-hub
Environment="BUNDLES_DIR=/var/lib/edge-policy-hub/config/tenants.d"
Environment="ENFORCER_HOST=127.0.0.1"
Environment="ENFORCER_PORT=8181"
Environment="LOG_LEVEL=info"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=edge-policy-enforcer
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/edge-policy-hub
CapabilityBoundingSet=
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
Alias=enforcer.service
