# infra/systemd/edge-policy-audit-store.service
[Unit]
Description=Edge Policy Hub - Audit Store Service
Documentation=https://github.com/acailic/edge-policy-hub
After=network-online.target
Wants=network-online.target
PartOf=edge-policy-hub.target

[Service]
Type=simple
ExecStart=/usr/local/bin/edge-policy-audit-store
Restart=on-failure
RestartSec=5s
User=edge-policy
Group=edge-policy
WorkingDirectory=/var/lib/edge-policy-hub
Environment="AUDIT_HOST=127.0.0.1"
Environment="AUDIT_PORT=8182"
Environment="AUDIT_DATA_DIR=/var/lib/edge-policy-hub/data/audit"
Environment="AUDIT_HMAC_SECRET_FILE=/etc/edge-policy-hub/hmac-secret"
Environment="ENABLE_DEFERRED_UPLOAD=true"
Environment="LOG_LEVEL=info"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=edge-policy-audit-store
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/edge-policy-hub/data/audit
ReadOnlyPaths=/etc/edge-policy-hub
CapabilityBoundingSet=
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
Alias=audit-store.service
